# templates/app.service.j2
{#
TODO: reuse your unit from Lab 1
TODO: Path to jar should use variables {{ app_dir }}/{{ app_file }}
TODO: Environment variable for PORT should use value from variable {{ app_port }}
#}
[Unit]
Description=Ansible_lab

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/bin/java -jar /opt/sample/service.jar
User=aleks
Environment=PORT="8080"
Restart=always
